# busylight.plist.in
<?xml version="1.0" encoding="UTF-8"?> 
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd"> 
<plist version="1.0"> 
    <dict> 
        <key>Label</key> 
            <string>com.engelschall.busylight</string> 
        <key>ProgramArguments</key> 
            <array> 
                <string>@basedir@/busylight.sh</string> 
            </array> 
        <key>WorkingDirectory</key> 
            <string>@basedir@</string> 
        <key>StandardOutPath</key> 
            <string>@basedir@/busylight.log</string> 
        <key>StandardErrorPath</key> 
            <string>@basedir@/busylight.log</string> 
        <key>StartInterval</key>
	        <integer>3600</integer>
        <key>KeepAlive</key> 
            <dict> 
                <key>SuccessfulExit</key> 
                <false/> 
            </dict> 
    </dict> 
</plist>
