# busylight.service.in
##
##  busylight -- Control Kuando BusyLight through REST
##  Copyright (c) 2025 Dr. Ralf S. Engelschall <rse@engelschall.com>
##  Licensed under MIT license <https://spdx.org/licenses/MIT.html>
##

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=@basedir@/busylight.sh

[Install]
WantedBy=multi-user.target

